# rootfs/etc/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysql/mysqld.pid /usr/sbin/mysqld --user mysql --pid-file=/var/run/mysql/mysqld.pid
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/data/mysql"